# config/supervisord.conf
[supervisord]
nodaemon=true

[program:umurmur]
command=/usr/bin/umurmurd -d -c /etc/umurmur/umurmur.conf -r
autorestart=true
